# chromeos-base/openfyde-utils/files/init/fydeos-collect-log.conf
# Copyright 2023 Fyde Innovations
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Collect fydeos log"
author          "fangzhou@fydeos.io"

start on stopped fydeos-scripts-mounter

oom score -100

exec /usr/share/cros/init/collect_fydeos_log.sh
