# chromeos-base/fydeos-hardware-tuner/files/fydeos-hwtuner-updater.conf
# Copyright (c) 2022 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Update fydeos hardware tuner files"
author          "fangzhou@fydeos.io"

oom score 0

start on shill-connected

script
	updater_bin="/usr/share/fydeos_shell/scripts_updater/updater.sh"
	tuner_dir="/usr/share/hwtuner-script"
	if [ -x "$updater_bin" ] && [ -d "$tuner_dir" ]; then
		tmp_log="/tmp/hwtuner-updater.log"
		"$updater_bin" -t "$tuner_dir" > "$tmp_log" 2>&1
	fi
end script
