# app-i18n/google-ime-tools/files/init/fydeos-ime-connection-checker.conf
# Copyright 2018 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Run the checker script for google-ime-tools"
author "fangzhou@fydeos.io"

start on started network-services

script
    exec /usr/share/cros/init/ime_connection_checker.sh
end script
